# backups/backup-home-nfs.service
[Unit]
Description=Backup /home to local nfs

[Service]
Type=oneshot
ExecStart=/root/scripts/backup-home-nfs.sh
