# backups/backup-home-nfs.service
[Unit]
Description=Backup /home to local nfs

[Service]
Type=simple
ExecStart=/root/scripts/backup-home-nfs.sh
